# home/dot_config/systemd/user/blueman-applet.service
[Unit]
Description=Blueman applet
PartOf=sway-session.target
After=sway-session.target

[Service]
ExecStart=/usr/bin/env blueman-applet
Restart=on-failure

[Install]
WantedBy=sway-session.target
