# home/dot_config/systemd/user/sway-session.target
[Unit]
Description=Sway session services
